# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; QUEUE_CONNECTION=database : sans ce worker, tout job mis en file (notifications,
; validations) s'empile dans la table `jobs` et n'est JAMAIS exécuté.
; --max-time relance le worker chaque heure : il libère la mémoire et recharge le
; code après un déploiement (un worker de longue durée garde l'ancien en mémoire).
[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --backoff=10 --max-time=3600
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
; Laisse au job en cours le temps de finir avant l'arrêt (--max-time=3600 + marge).
stopwaitsecs=65
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
